# debian/cron.d
# Every minute, do mailman actions
* * * * *	   root		/usr/lib/alternc/update_mailman.sh
